# systemd/cloudplay-dev.container
[Unit]
Description=Cloudplay dev container (long-lived iteration harness)
Wants=network-online.target
After=network-online.target

[Container]
ContainerName=cloudplay-dev
Image=localhost/cloudplay-dev:latest

# Match production sandboxing so GL/NVENC behavior is identical.
SecurityLabelDisable=true
# --security-opt seccomp=unconfined: libretro cores + NVIDIA userspace.
# --group-add keep-groups: propagates rocks's supplementary groups (notably
# the host's `input` gid 104) into the container namespace so /dev/uinput
# under mode 0660 root:input is writable. Without this, the container's
# GroupAdd=input creates a *separate* gid 996 that doesn't match the host.
PodmanArgs=--security-opt seccomp=unconfined --group-add keep-groups

# GPU + uinput passthrough.
AddDevice=/dev/dri
AddDevice=nvidia.com/gpu=all
AddDevice=/dev/uinput
# Bind-mount the whole /dev/input tree so the event nodes created by our
# VirtualPad (via UI_DEV_CREATE) are visible inside the container. AddDevice
# only exposes specific nodes; we can't predict which eventN udev will assign.
Volume=/dev/input:/dev/input

# Source tree and build-output lanes (rsync target + go build output).
Volume=/home/rocks/containers/cloudplay-dev/src:/src:z
Volume=/home/rocks/containers/cloudplay-dev/out:/out:z
# xemu BIOS + HDD. Must be rw because the Xbox HDD image is a qcow2 that
# xemu writes back to for save games, dashboard settings, eeprom, etc.
# BIOS roms under /xemu-bios/bios and /xemu-bios/mcpx are effectively static
# but we keep the whole dir on one mount for simplicity.
Volume=/home/rocks/containers/cloudplay-dev/xemu-bios:/xemu-bios:z
# Optional: test ROMs for harness use.
# games: rw so extract-xiso can rewrite in place and xemu can write save data
Volume=/home/rocks/containers/cloudplay-dev/games:/games:z

# Mirror production's Vulkan ICD + display env so behavior matches.
Volume=/usr/share/vulkan/icd.d:/usr/share/vulkan/icd.d:ro
Environment=DISPLAY=:99
Environment=MESA_GL_VERSION_OVERRIDE=4.5
Environment=NVIDIA_DRIVER_CAPABILITIES=all
Environment=LD_LIBRARY_PATH=/lib64:/usr/lib64:/usr/lib:/lib
Environment=__GLX_VENDOR_LIBRARY_NAME=nvidia
Environment=GOCACHE=/go/cache

# No PublishPort — we drive the container exclusively via `podman exec`.
# If you want to run the full worker here against a browser, uncomment and
# pick a distinct port from production (prod uses 8009/tcp + 9753/udp).
#PublishPort=0.0.0.0:8010:8000/tcp
#PublishPort=0.0.0.0:9754:9754/udp

[Service]
# Don't restart-loop a busted image; we'll fix and `systemctl --user restart`.
Restart=on-failure
TimeoutStopSec=30

[Install]
WantedBy=default.target
